# misc/maintenance.service
[Unit]
Description=Run maintenance script
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/octopus/maintenance_script.sh