# misc/maintenance.timer
[Unit]
Description=Run maintenance script timer

[Timer]
OnUnitActiveSec=24h
Unit=maintenance.service

[Install]
WantedBy=timers.target
